# DEBIAN/postinst
#!/bin/bash
systemctl enable dht22-sensor