# DEBIAN/prerm
#!/bin/bash
systemctl stop dht22-sensor
systemctl disable dht22-sensor
